# src/idrive_backup.timer
[Unit]
Description=IDrive backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
